# dump.service
[Unit]
Description=Dump database

[Service]
Type=simple
ExecStart=/path/to/file.sh
